# ql-ol-sdk/ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/busybox-mdev.list
/etc/mdev/usb.sh
/etc/mdev/iio.sh
/etc/mdev/
/etc/mdev/automountsdcard.sh
/etc/mdev/mdev-mount.sh
/etc/mdev/find-touchscreen.sh
/etc/mdev.conf
/etc/init.d/mdev
