# ql-ol-sdk/ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/diag.list
/etc/
/usr/bin/test_diag
/usr/bin/diag_socket_log
/usr/bin/diag_callback_sample
/usr/bin/diag_dci_sample
/usr/lib/libdiag.so.1.0.0
/etc/init.d/
/usr/bin/diag_uart_log
/etc/init.d/chgrp-diag
/usr/bin/PktRspTest
/usr/lib/libdiag.so.1
/usr/bin/diag_mdlog
/usr/bin/diag_klog
