# ql-ol-sdk/ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/miniupnpd.list
/data/miniupnpd/miniupnpd.conf
/usr/sbin/miniupnpd
/data/miniupnpd/
/etc/init.d/miniupnpd
