# ql-ol-sdk/ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/minidlna.list
/usr/sbin/minidlnad
/etc/init.d/minidlna
/data/minidlna.conf
