# ql-ol-sdk/ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/ql-forward.list
/etc/init.d/start_ql_forward_le
/usr/bin/ql_forward
