# ql-ol-sdk/ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/qmi-framework.postinst
if true; then
      [ -n "$D" ] && OPT="-r $D" || OPT="-s"
      update-rc.d $OPT -f init_irsc_util remove
      update-rc.d $OPT init_irsc_util start 29 2 3 4 5 . stop 71 0 1 6 .
     fi
if type update-rc.d >/dev/null 2>/dev/null; then
	if [ -n "$D" ]; then
		OPT="-r $D"
	else
		OPT="-s"
	fi
	update-rc.d $OPT init_irsc_util start 29 2 3 4 5 . stop 71 0 1 6 .
fi
