# ql-ol-sdk/ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/quectel-gps-handle.list
/usr/bin/quectel-gps-handle
/usr/bin/quectel-uart-ddp
/etc/init.d/quectel-gps-handle
