# ql-ol-sdk/ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/quectel-psm-aware.list
/usr/bin/quectel_psm_aware
/etc/init.d/quectel_psm_aware
