# ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/atfwd-daemon.list
/usr/bin/wget
/usr/bin/atfwd_daemon
/etc/init.d/start_atfwd_daemon
