# ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/avahi-dnsconfd.list
/usr/sbin/avahi-dnsconfd
/etc/init.d/avahi-dnsconfd
/etc/avahi/avahi-dnsconfd.action
