# ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/data.list
/usr/lib/libqdp.so.0.0.0
/etc/init.d/start_QCMAP_ConnectionManager_le
/WEBSERVER/www/QCMAP_Firewall.html
/etc/udhcpc.d/
/WEBSERVER/www/QCMAP_Firewall_Help.html
/usr/lib/libqdi.so.0.0.0
/usr/lib/libqcmapipc.so.1.0.0
/etc/portal-gateway.conf
/WEBSERVER/www/QCMAP_IPV6_Help.html
/WEBSERVER/www/QCMAP_login.html
/WEBSERVER/
/usr/bin/qmi_ip_multiclient
/usr/bin/portal_ctl
/usr/tests/dsi_netctrl_test
/usr/lib/libqcmap_cm.so.1
/WEBSERVER/www/QCMAP_WLAN.html
/data/mobileap_cfg.xml
/usr/bin/netmgrd
/etc/init.d/port_bridge
/etc/init.d/start_qti_le
/WEBSERVER/www/images/delete.png
/usr/lib/libqcmapipc.so.1
/usr/bin/radish
/usr/bin/QCMAP_CLI
/etc/init.d/start_stop_qti_ppp_le
/etc/data/
/usr/lib/libqcmaputils.so.1.0.0
/WEBSERVER/www/QCMAP_UPNP_DLNA_MDNS.html
/usr/bin/eMBMs_TunnelingModule
/WEBSERVER/www/cgi-bin/qcmap_web_cgi
/usr/lib/libqmi_ip.so.1.0.0
/usr/lib/libqdi.so.0
/usr/lib/libqcmap_cm.so.1.0.0
/WEBSERVER/www/QCMAP_Help.html
/WEBSERVER/www/QCMAP_Account_Help.html
/WEBSERVER/www/QCMAP_IPV4.html
/usr/lib/libqmi_ip.so.1
/usr/bin/portal_gateway
/etc/data/netmgr_config.xml
/data/qmi_ip_cfg.xml
/WEBSERVER/www/QCMAP_IPV4_Help.html
/usr/lib/libqdp.so.0
/usr/tests/
/etc/init.d/start_eMBMs_TunnelingModule_le
/WEBSERVER/www/images/
/etc/qmi_ip_cfg.xml
/usr/lib/libnetmgr.so.0.0.0
/WEBSERVER/www/QCMAP_WLAN_Help.html
/WEBSERVER/www/QCMAP_WWAN.html
/WEBSERVER/www/QCMAP.html
/etc/init.d/netmgrd
/etc/init.d/data-init
/etc/udhcpc.d/udhcpc.script
/usr/bin/QCMAP_Web_CLIENT
/usr/lib/libdsi_netctrl.so.0
/usr/bin/QCMAP_ConnectionManager
/WEBSERVER/www/QCMAP_WWAN_Help.html
/usr/bin/qti
/WEBSERVER/www/cgi-bin/
/WEBSERVER/www/QCMAP_UPNP_DLNA_MDNS_Help.html
/WEBSERVER/www/QCMAP_IPV6.html
/WEBSERVER/www/images/header.jpg
/usr/bin/QCMAP_StaInterface
/WEBSERVER/www/QCMAP_Account.html
/usr/lib/libqcmap_client.so.1
/WEBSERVER/www/QCMAP_NAT.html
/WEBSERVER/www/
/usr/lib/libqcmap_client.so.1.0.0
/WEBSERVER/www/QCMAP_NAT_Help.html
/usr/lib/libqcmaputils.so.1
/data/
/usr/bin/qti_ppp
/etc/portal-gateway.conf.sample
/etc/init.d/start_stop_qmi_ip_multiclient
/WEBSERVER/www/cgi-bin/qcmap_auth
/usr/lib/libdsi_netctrl.so.0.0.0
/WEBSERVER/www/images/update.jpg
/usr/lib/libnetmgr.so.0
/etc/data/dsi_config.xml
