# ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/dnsmasq.list
/etc/dnsmasq.d/
/data/dnsmasq.conf
/usr/bin/dnsmasq
/etc/init.d/dnsmasq
/bin/dnsmasq_script.sh
