# ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/fs-scrub-daemon.list
/etc/init.d/start_fs-scrub-daemon
/sbin/fs-scrub-daemon
