# ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/loc-launcher.list
/usr/bin/loc_launcher
/etc/init.d/start_loc_launcher
