# ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/mcm-core.list
/usr/bin/MCM_MOBILEAP_ConnectionManager
/usr/bin/mcm_data_srv
/etc/init.d/start_MCM_atcop_svc_le
/etc/init.d/start_MCM_MOBILEAP_ConnectionManager_le
/usr/bin/MCM_MOBILEAP_CLI
/usr/lib/libmcm.so.0.0.0
/usr/bin/MCM_atcop_svc
/usr/lib/libmcm_log_util.so.0
/usr/bin/uim_test_client
/etc/init.d/start_mcm_data_srv_le
/usr/lib/libmcm_srv_audio.so.0.0.0
/usr/bin/mcm_ril_service
/usr/lib/libmcmipc.so.0
/usr/lib/libmcm_srv_audio.so.0
/usr/tests/mcm_data_test
/usr/lib/libmcmipc.so.0.0.0
/etc/init.d/start_mcm_ril_srv_le
/usr/bin/wlan_test_client
/usr/lib/libmcm.so.0
/usr/bin/MCM_ATCOP_CLI
/usr/bin/IoEConsoleClient
/usr/lib/libmcm_log_util.so.0.0.0
