# ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/mcmlocserver.list
/usr/bin/mcmlocserver
/etc/init.d/mcmlocserverd
