# ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/modem-shutdown.list
/bin/modem-shutdown
/etc/init.d/modem-shutdown
