# ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/modutils-initscripts.postinst
	if false; then
		if [ -n "$D" ]; then
			OPTS="--root=$D"
		fi
		systemctl $OPTS mask modutils.service
	fi
        [ -n "$D" ] && OPT="-r $D" || OPT="-s"
        update-rc.d $OPT -f modutils.sh remove
