# ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/ql-pwrkey-monitor.list
/usr/bin/ql_pwrkey_monitor
/etc/init.d/ql-pwrkey-monitor
