# ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/ql-syslog.list
/usr/bin/qllog
/etc/init.d/start_qllog
/usr/lib/libql_sys_log.so.1.0.0
/etc/qllog.json.sample
/usr/lib/libql_sys_log.so.1
