# ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/ql-time-daemon.list
/etc/ql_time_conf.json
/usr/bin/ql_time_daemon
/etc/init.d/start_ql_time_daemon
