# ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/qmuxd.list
/etc/data/qmi_config.xml
/etc/init.d/qmuxd
/usr/bin/qmuxd
