# ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/quectel-pcm-daemon.list
/etc/init.d/aaa_quectel_pcm_daemon
/usr/bin/quectel_pcm_daemon
