# ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/quectel-daemon.list
/usr/bin/quectel_daemon
/etc/ringtone1.wav
/etc/init.d/quectel_daemon
