# ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/quectel-smd-atcmd.list
/etc/init.d/quectel-smd-atcmd
/usr/bin/quectel-smd-atcmd
