# ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/quectel-monitor-daemon.list
/usr/bin/quectel_monitor_daemon
/etc/init.d/quectel-monitor-daemon
