# ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/quectel-remotefs.list
/usr/bin/quectel-remotefs-service
/etc/init.d/start_quectel_remotefs_daemon
