# ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/quectel-thermal.list
/usr/bin/quectel-thermal
/etc/init.d/quectel-thermal
