# ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/subsystem-ramdump.list
/etc/init.d/start_subsystem_ramdump
/usr/bin/subsystem_ramdump
