# ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/tftp-server.list
/system/rfs/
/sbin/tftp_server
/data/persist/rfs/
/system/rfs/msm/mpss/readonly/
/system/rfs/mdm/mpss/readonly/firmware
/data/persist/hlos_rfs/shared/
/system/rfs/msm/mpss/readwrite
/data/tombstones/lpass/
/data/persist/rfs/mdm/mpss/
/system/rfs/mdm/adsp/readonly/
/data/persist/rfs/msm/
/system/rfs/mdm/mpss/
/data/persist/
/data/persist/rfs/shared/
/data/persist/rfs/mdm/
/system/rfs/msm/adsp/readonly/firmware
/data/persist/hlos_rfs/
/system/rfs/mdm/mpss/ramdumps
/system/rfs/mdm/adsp/readwrite
/system/rfs/msm/mpss/ramdumps
/system/rfs/msm/mpss/readonly/firmware
/system/rfs/mdm/mpss/shared
/data/tombstones/modem/
/data/persist/rfs/mdm/adsp/
/system/rfs/msm/
/data/persist/rfs/msm/mpss/
/system/rfs/mdm/adsp/
/system/rfs/msm/mpss/shared
/system/rfs/mdm/
/system/rfs/mdm/mpss/readonly/
/system/rfs/msm/mpss/
/system/rfs/mdm/adsp/shared
/system/rfs/mdm/mpss/readwrite
/system/rfs/mdm/adsp/hlos
/system/rfs/msm/adsp/shared
/system/rfs/msm/adsp/readonly/
/data/persist/rfs/msm/adsp/
/system/rfs/mdm/mpss/hlos
/system/rfs/mdm/adsp/ramdumps
/system/rfs/msm/adsp/
/system/rfs/msm/adsp/readwrite
/system/rfs/msm/adsp/hlos
/system/rfs/msm/adsp/ramdumps
/system/rfs/msm/mpss/hlos
/data/tombstones/
/system/
/system/rfs/mdm/adsp/readonly/firmware
/etc/init.d/tftp_server.sh
