# ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/time-services.list
/usr/bin/time_daemon
/etc/init.d/time_serviced
