# ql-ol-crosstool/sysroots/armv7a-vfp-neon-oe-linux-gnueabi/var/lib/opkg/info/wlan-services.preinst
#!/bin/sh
if [ -z "$D" -a -f "/etc/init.d/start_wlan_services" ]; then
	/etc/init.d/start_wlan_services stop
fi
if type update-rc.d >/dev/null 2>/dev/null; then
	if [ -n "$D" ]; then
		OPT="-f -r $D"
	else
		OPT="-f"
	fi
	update-rc.d $OPT start_wlan_services remove
fi
